# scripts/ratan-engine.service
[Unit]
Description=RATAN Engine - Resilient Aggregation and Transition for Advanced Networking
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ratan-engine
EnvironmentFile=-/etc/ratan-engine/env
ExecStart=/usr/bin/python3 main.py --mode server
Restart=always
RestartSec=5
Environment=RATAN_PROJECT_DIR=/root/ratan-engine
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
